# services/minecraft-backup.service
[Unit]
Description=Daily folder backup

[Service]
Type=oneshot
ExecStart=/home/npi5/backup-script/minecraft-backup.sh

